# config/supervisor/uwsgi.conf
[program:uwsgi]
command=/var/deploy/fitting/env/bin/uwsgi --virtualenv /var/deploy/fitting/env/ -i /var/wsgi/fitting/config/uwsgi.ini
redirect_stderr=true
directory=/var/wsgi/fitting/
